# habit-api.service
[Unit]
Description=Habit API
After=network.target

[Service]
Type=simple
User=agentp
WorkingDirectory=/home/agentp/habit-api
ExecStart=/home/agentp/.nvm/versions/node/v20.20.2/bin/node index.js
Restart=always
RestartSec=5
StandardOutput=append:/home/agentp/habit-api/app.log
StandardError=append:/home/agentp/habit-api/app.log

[Install]
WantedBy=multi-user.target
